# systemd/makeros-hub.service
[Unit]
Description=makeros native-print hub agent
Documentation=https://github.com/overengineered-solutions/makeros-hub-agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=makeros-hub
Group=makeros-hub
Environment=PYTHONPATH=/opt/makeros-hub
# Virtual-printer FTPS binds 990 without running the whole service as root.
AmbientCapabilities=CAP_NET_BIND_SERVICE
# The venv python carries paho-mqtt (the Bambu adapter). install.sh creates it.
# If it's ever absent the agent still enrolls + heartbeats — swap to
# /usr/bin/python3 and printers report 'agent_missing_paho' until paho is added.
ExecStart=/opt/makeros-hub/venv/bin/python -m makeros_hub run
# Restart=always IS the self-healing reconnect primitive — a network blip,
# a cloud deploy, or a transient 5xx just bounces the loop.
Restart=always
RestartSec=5
TimeoutStopSec=25
# Modest hardening (a box on the shop LAN). The credential + state live under
# the StateDirectory, which systemd creates 0700 owned by the service user.
# NOTE: deliberately NO NoNewPrivileges — it blocks the agent's one privileged
# path, the sudo-scoped OTA self-update (`sudo /opt/makeros-hub/update.sh`),
# which is the ONLY thing the service user may run as root. With NoNewPrivileges
# set, sudo can't escalate even with a valid NOPASSWD rule, so OTA silently
# fails ("no new privileges flag is set"). The narrow sudoers rule is the real
# control here.
# NOTE: deliberately NO ProtectSystem — the agent legitimately installs/manages
# system software (Tailscale via apt, and OTA self-update), which write to /usr
# (e.g. /usr/share/keyrings, /usr/bin) and /etc. ProtectSystem=full mounts those
# read-only, so even the sudo-scoped install fails ("Read-only file system" when
# tee'ing the Tailscale APT keyring) — confirmed live via the diagnostics bundle.
# The unprivileged service user + the narrow sudoers rules (update.sh /
# tailscale-setup.sh only) are the real control. ProtectHome + PrivateTmp stay.
ProtectHome=true
PrivateTmp=true
StateDirectory=makeros-hub
ConfigurationDirectory=makeros-hub

[Install]
WantedBy=multi-user.target
