# systemd/makeros-hub.service
[Unit]
Description=makeros native-print hub agent
Documentation=https://github.com/overengineered-solutions/makeros-hub-agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=makeros-hub
Group=makeros-hub
Environment=PYTHONPATH=/opt/makeros-hub
# The venv python carries paho-mqtt (the Bambu adapter). install.sh creates it.
# If it's ever absent the agent still enrolls + heartbeats — swap to
# /usr/bin/python3 and printers report 'agent_missing_paho' until paho is added.
ExecStart=/opt/makeros-hub/venv/bin/python -m makeros_hub run
# Restart=always IS the self-healing reconnect primitive — a network blip,
# a cloud deploy, or a transient 5xx just bounces the loop.
Restart=always
RestartSec=5
# Modest hardening (a box on the shop LAN). The credential + state live under
# the StateDirectory, which systemd creates 0700 owned by the service user.
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
StateDirectory=makeros-hub
ConfigurationDirectory=makeros-hub

[Install]
WantedBy=multi-user.target
